# config/cron/api.crontab
# intentionally empty, no crontasks run on the production machine that hosts the API frontend and mongodb
0 6 * * *   cd /projects/rtc/current && rake analytics:report