# config/cron/api.crontab
0 6 * * *  /projects/rtc/shared/cron/analytics.sh
