# config/cron/api.crontab
0 6 * * *   cd /projects/rtc/current && rake analytics:report
