# config/cron/api.crontab
0 6 * * *  /projects/congress/shared/cron/analytics.sh

# leave newline
