# config/cron/staging.crontab
# legislators and committees from Sunlight API
0 3 * * *    /projects/congress/shared/cron/task.sh legislators && /projects/congress/shared/cron/task.sh committees

# bills from THOMAS, amendments through GovTrack
0 4,16 * * * /projects/congress/shared/cron/task.sh bills_thomas && /projects/congress/shared/cron/task.sh amendments_archive

# full bill text from GPO
*/120 * * * * /projects/congress/shared/cron/task.sh bulk_gpo_bills
0 5 * * * /projects/congress/shared/cron/task.sh bill_text_archive

# popular bill names, check once a day (this is staging)
3 * * * * /projects/congress/shared/cron/task.sh bills_popular_names

# video
0 9 * * *    /projects/congress/shared/cron/task.sh videos captions=True
0 7 * * *    /projects/congress/shared/cron/task.sh videos captions=True senate=True


# votes -

# download all house votes for the last 10 votes every 2 hours (staging)
# don't bother re-downloading votes for the year on staging
*/120 * * * * /projects/congress/shared/cron/task.sh votes_house latest=10
*/120 * * * * /projects/congress/shared/cron/task.sh votes_senate latest=10


# floor updates
*/120 * * * *  /projects/congress/shared/cron/task.sh floor_senate
*/120 * * * *  /projects/congress/shared/cron/task.sh floor_house

# committee hearings
*/120 * * * * /projects/congress/shared/cron/task.sh hearings_senate
*/120 * * * * /projects/congress/shared/cron/task.sh hearings_house

# GAO reports (full text)
10 * * * * /projects/congress/shared/cron/task.sh documents_gao_reports cache=1

# upcoming leadership information
*/120 * * * * /projects/congress/shared/cron/task.sh upcoming_senate
*/120 * * * * /projects/congress/shared/cron/task.sh upcoming_house

# Regulations
*/119 * * * * /projects/congress/shared/cron/task.sh regulations_archive cache=1
*/120 * * * * /projects/congress/shared/cron/task.sh regulations_archive cache=1 public_inspection=1


# leave the newline here
